# raspberrypi/Pi-hole/backup.cron
0 * * * * /home/pi/tools/backup-pihole.sh >/dev/null 2>&1